# dualboot_android/important files after compiling/data/etc/init/audiod.conf
start on (started udev
             and started dbus) 

script
export mpt=/webos
PATH=$mpt/usr/sbin:/$mpt/usr/bin:$mpt/sbin:$mpt/bin:$PATH
LD_LIBRARY_PATH=$mpt/usr/lib:$mpt/lib \
$mpt/lib/ld-linux.so.3 $mpt/usr/sbin/audiod -d $* > /var/log/audiod.log || exit 0
end script
