# important files after compiling/data/etc/init/killudev.conf
start on starting ubuntu-touch-session

script
pkill "systemd-udevd"
end script

