# important files after compiling/data/etc/init/killudev3.conf
start on (local-filesystems
          and started dbus
          and static-network-up)

script
pkill "systemd-udevd"
end script

