# important files after compiling/data/etc/init/killudev1.conf
start on mounted

script
kill 125
end script

