# important files after compiling/data/etc/init/killudev2.conf
start on (startup
          and starting udevmonitor)

script
pkill "systemd-udevd"
end script

