# important files after compiling/data/etc/init/powerd.conf
author "Michael Frey <michael.frey@canonical.com>"

start on started dbus and android

respawn

env ANDROID_ROOT=/system
#uncomment the line below to enable debugging
env POWERD_DEBUG=1

exec /usr/bin/powerd
