# important files after compiling/data/etc/init/ubuntu-touch-session-setup.conf
author "Ricardo Mendoza <ricmm@ubuntu.com>"

start on (started udev
	     and started dbus)
stop on stopping dbus

env user=phablet
export user

script
    sleep 30
    start ubuntu-touch-session USER=phablet
end script
